# launchd/com.example.vibeproxy-canary.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template. install.sh substitutes __HOME__ and __PYTHON__.
  __PYTHON__ must have PyYAML (canary imports yaml). install.sh verifies.
  Probes the Claude lane every 6h; alerts via CANARY_NOTIFY_CMD on persistent failure.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.example.vibeproxy-canary</string>
	<key>ProgramArguments</key>
	<array>
		<string>__PYTHON__</string>
		<string>__HOME__/.local/share/vibeproxy-claude-lane/scripts/canary.py</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CANARY_SHIM_URL</key>
		<string>http://127.0.0.1:8485</string>
		<key>CANARY_UPSTREAM_URL</key>
		<string>http://127.0.0.1:8318</string>
		<key>CANARY_MODEL</key>
		<string>claude-sonnet-5</string>
		<key>CANARY_MERGED_CONFIG</key>
		<string>__HOME__/.cli-proxy-api/merged-config.yaml</string>
		<key>CANARY_MERGER</key>
		<string>__HOME__/.local/share/vibeproxy-claude-lane/scripts/merge-config-extras.py</string>
		<key>CANARY_NOTIFY_CMD</key>
		<string>hermes send -t telegram -q</string>
		<key>CANARY_LOG</key>
		<string>__HOME__/.local/state/vibeproxy-claude-lane/canary.log</string>
	</dict>
	<key>StartInterval</key>
	<integer>21600</integer>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardOutPath</key>
	<string>__HOME__/.local/state/vibeproxy-claude-lane/canary.launchd.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/.local/state/vibeproxy-claude-lane/canary.launchd.log</string>
</dict>
</plist>
